# archup.service
[Unit]
Description=Script to check for available package updates

[Service]
ExecStart=/usr/bin/archup -c

[Install]
WantedBy=multi-user.target
